# files/etc/uci-defaults/99_network.erb
#!/bin/sh
cat /usr/local/lib/mesh.uci | uci -m import network
cat /usr/local/lib/allow_ssh.uci | uci -m import firewall
uci set network.loopback.ip6addr='fda1:987a:cc9:1::1/64'
uci set network.lan.ip6addr='fda1:987a:cc9::1/48'

uci set firewall.@zone[0].network='lan mesh'

uci commit network
uci commit firewall


