# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

iperf -s -D > /tmp/iperf_tcp.log
iperf -s -u -D > /tmp/iperf_udp.log
/usr/local/bin/dump_status.sh &
/usr/local/bin/commit_dumps.sh &

exit 0
